# startcask/systemd/plasma-ksplash-ready.service.in
[Unit]
Description=KSplash "ready" Stage
Wants=plasma-core.target
After=plasma-core.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=-@QtBinariesDir@/qdbus org.kde.KSplash /KSplash org.kde.KSplash.setStage ready
Slice=session.slice
