# startcask/systemd/plasma-workspace-x11.target
[Unit]
Wants=plasma-kwin_x11.service
Requires=plasma-workspace.target
BindsTo=plasma-ksmserver.service
