# startcask/systemd/plasma-workspace-wayland.target
[Unit]
Requires=plasma-workspace.target
Requires=plasma-kwin_wayland.service
BindsTo=plasma-kwin_wayland.service
